# assets/etc/systemd/system/docker.service.d/custom-exec-start.conf.template
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon -H fd:// --iptables=false{{#each host.variables.nameservers}} --dns {{{ . }}}{{/each}}
